# models/hermes/systemd/kortexa-ai-llm-hermes.service
[Unit]
Description=Hermes LLM
After=network.target

StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple

User=francip
Group=francip

WorkingDirectory=/home/francip/src/llm.server/models/hermes

Environment=PATH=/home/francip/src/llm.server/models/hermes:/home/francip/src/llm.server/models/hermes/.venv/bin
Environment=PYTHONUNBUFFERED=1

ExecStart=/home/francip/src/llm.server/models/hermes/run.sh

Restart=always
RestartSec=15

TimeoutStopSec=60

ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/francip/src/llm.server /home/francip/.cache /home/francip/.config /home/francip/data/kortexa-ai
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
